# src/app/crontab
0 2 * * 0 python3 src/app/spotifyautoshuffle.py >> /logfile