# src/app/crontab
0 2 * * 0 python /spotifyautoshuffle.py >> /logfile